# deployment/systemd/pow-service.service
[Unit]
Description=Rakuten POW Service
Documentation=https://github.com/your-org/rakuten-checker
After=network.target redis.service
Wants=redis.service
Requires=network.target

[Service]
Type=simple
User=rakuten
Group=rakuten
WorkingDirectory=/opt/rakuten-checker
ExecStart=/usr/bin/node src/pow-service/index.js
EnvironmentFile=/opt/rakuten-checker/.env

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Resource limits (POW service is CPU intensive)
LimitNOFILE=65536
LimitNPROC=8192
CPUQuota=800%

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/rakuten-checker/logs

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pow-service

# Graceful shutdown
TimeoutStopSec=60
KillMode=mixed
KillSignal=SIGTERM

# Environment
Environment=NODE_ENV=production
Environment=LOG_LEVEL=info
Environment=PORT=3001

[Install]
WantedBy=multi-user.target
