# deployment/pow-service.service
[Unit]
Description=Rakuten POW Service
Documentation=https://github.com/your-org/rakuten-checker
After=network.target redis.service
Wants=redis.service
Requires=network.target

[Service]
Type=simple
User=pow-service
Group=pow-service
WorkingDirectory=/opt/rakuten-checker
ExecStart=/usr/bin/docker run --rm --name pow-service \
  --network host \
  --env-file /opt/rakuten-checker/.env.pow-service \
  --volume /opt/rakuten-checker/logs:/app/logs \
  rakuten-pow-service:latest

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Resource limits (POW service is CPU intensive)
LimitNOFILE=65536
LimitNPROC=8192
CPUQuota=800%  # Allow up to 8 CPU cores

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/rakuten-checker/logs

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pow-service

# Graceful shutdown
TimeoutStopSec=60
KillMode=mixed
KillSignal=SIGTERM

# Environment
Environment=NODE_ENV=production
Environment=LOG_LEVEL=info
Environment=PORT=8080

[Install]
WantedBy=multi-user.target