# deployment/worker.service
[Unit]
Description=Rakuten Worker Service
Documentation=https://github.com/your-org/rakuten-checker
After=network.target redis.service coordinator.service
Wants=redis.service coordinator.service
Requires=network.target

[Service]
Type=simple
User=worker
Group=worker
WorkingDirectory=/opt/rakuten-checker
ExecStart=/usr/bin/docker run --rm --name worker-%i \
  --network host \
  --env-file /opt/rakuten-checker/.env.worker \
  --env WORKER_ID=worker-%i \
  --volume /opt/rakuten-checker/logs:/app/logs \
  rakuten-worker:latest

# Restart policy
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=5

# Resource limits
LimitNOFILE=32768
LimitNPROC=2048

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/rakuten-checker/logs

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=worker-%i

# Graceful shutdown
TimeoutStopSec=120
KillMode=mixed
KillSignal=SIGTERM

# Environment
Environment=NODE_ENV=production
Environment=LOG_LEVEL=info
Environment=WORKER_CONCURRENCY=3

[Install]
WantedBy=multi-user.target